# roles/install_bird_lgproxy_go/templates/bird-lgproxy.service.j2
[Unit]
Description=bird-lg proxy

[Service]
User=bird
Group=bird
ExecStart=/usr/local/bin/bird-lgproxy-go

[Install]
WantedBy=multi-user.target
