# roles/install_pingfinder/templates/pingfinder.timer.j2
[Unit]
Description=Run DN42 Pingfinder periodically

[Timer]
OnBootSec=3m
OnUnitActiveSec=3m
AccuracySec=1m

[Install]
WantedBy=timers.target
