# roles/install_pingfinder/templates/pingfinder.timer.j2
[Unit]
Description=Run DN42 Pingfinder periodically

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m
AccuracySec=1m

[Install]
WantedBy=timers.target
