# roles/install_pingfinder/templates/pingfinder.service.j2
[Unit]
Description=Run DN42 Pingfinder

[Service]
Type=oneshot
WorkingDirectory=/opt/pingfinder
ExecStart=/opt/pingfinder/pingfinder.sh
