# roles/install_bird_lg_go/templates/bird-lg.service.j2
[Unit]
Description=bird-lg

[Service]
User=bird
Group=bird
ExecStart=/usr/local/bin/bird-lg-go
WorkingDirectory=/usr/local/bin

[Install]
WantedBy=multi-user.target
